# app/app.service
[Unit]
Description=Gunicorn instance to serve app
After=network.target
[Service]
User=ktoto
Group=lefttech